# root/etc/services.d/domoticz/run
#!/usr/bin/with-contenv bash
exec \
	s6-setuidgid abc /var/lib/domoticz/domoticz \
		-sslwww 1443 \
		-sslcert /config/keys/server_cert.pem \
		-userdata /config/ \
		-dbase /config/domoticz.db
